# deploy/baota/supervisor.ini
; 宝塔 → 软件商店 → 安装「Supervisor管理器」→ 添加守护进程
; 或使用「Python项目管理器」等效配置
;
; 名称: diagram-api
; 运行目录: /www/wwwroot
; 启动命令见下方 command 行

[program:diagram-api]
directory=/www/wwwroot
command=/www/wwwroot/venv/bin/python -m uvicorn backend.main:app --host 127.0.0.1 --port 8765 --workers 1
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
user=www
redirect_stderr=true
stdout_logfile=/www/wwwlogs/diagram-api.log
stdout_logfile_maxbytes=20MB
environment=PATH="/www/wwwroot/venv/bin",PLANTUML_SERVER="https://www.plantuml.com/plantuml/svg/"
